# rootfs/etc/systemd/system/mullvad-daemon.service.d/override.conf
[Service]
Environment="TALPID_DISABLE_OFFLINE_MONITOR=1"