# rootfs/etc/systemd/system/mullvad-suspend.service
[Unit]
Before=suspend.target

[Service]
Type=simple
ExecStart=/usr/bin/systemctl stop mullvad-daemon.service

[Install]
WantedBy=suspend.target
